# arkanelinux-generic/overlay/post_install/etc/udev/rules.d/83-trackpad.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="0265", RUN+="/usr/bin/bash -c 'modprobe -r hid_magicmouse'"
